# roles/node/templates/node.service.j2
[Service]
WorkingDirectory=/var/www/{{inventory_hostname}}
